# .config/systemd/user/lock.target
[Unit]
Description=User session lock target
Documentation=man:systemd.special(7)
RefuseManualStart=yes
StopWhenUnneeded=yes
