# .config/systemd/user/waypaper.service
[Unit]
Description=Restore wallpaper with Waypaper
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=oneshot
ExecStart=/usr/bin/waypaper --restore
RemainAfterExit=yes

[Install]
WantedBy=graphical-session.target
